# SD-Card-Contents/home/buccaneer/bin/buccaneer-wifi/2.1/configFiles/wifiDongle-detect.rules
ACTION=="add", ATTRS{idVendor}=="1d6b", RUN+="/home/buccaneer/bin/buccaneer-wifi/current/staticWifiNames.sh"
ACTION=="remove", ATTRS{idVendor}=="1d6b", RUN+="/home/buccaneer/bin/buccaneer-wifi/current/staticWifiNames.sh"
